# conf_server/systemd/system/backup-tarsnap.timer
[Unit]
Description=Run backup-tarsnap.service daily

[Timer]
OnCalendar=*-*-* 6:10:00
