# conf_server/systemd/system/backup-tarnsap.timer
[Unit]
Description=Run backup-tarnsap.service daily

[Timer]
OnCalendar=*-*-* 6:10:00
