# conf_server/systemd/system/backup-raid.service
# requires:
# rsync_time_machine	https://github.com/laurent22/rsync-time-backup

[Unit]
Description=Backup to RAID array

[Service]
Type=oneshot
ExecStart=/bin/sh -c ' \
	rsync_tm \
		/opt/containers/ \
		/mnt/media/backup/containers \
		/opt/containers/.backup-ignore 2>&1 \
	| mail -s "status for backup-raid.service" -- senan@senan.xyz \
'
