# conf_server/systemd/system/backup-raid.timer
[Unit]
Description=Run backup-raid.service daily

[Timer]
OnCalendar=*-*-* 6:00:00

[Install]
WantedBy=basic.target
