# conf_server/systemd/system/clean-docker.timer
[Unit]
Description=Run clean-docker.service daily

[Timer]
OnCalendar=*-*-* 5:50:00

[Install]
WantedBy=basic.target
