# conf_server/systemd/system/clean-docker.service
[Unit]
Description=Docker system prune

[Service]
Type=oneshot
ExecStart=/usr/bin/docker system prune -f
