# conf_server/systemd/system/clean-jellyfin.timer
[Unit]
Description=Run clean-jellyfin.service weekly

[Timer]
OnCalendar=weekly

[Install]
WantedBy=basic.target
