# conf_server/systemd/system/clean-gonic.timer
[Unit]
Description=Run clean-gonic.service weekly

[Timer]
OnCalendar=weekly

[Install]
WantedBy=basic.target
