# conf_server/systemd/system/email-tree.timer
[Unit]
Description=Run email-tree.service weekly 

[Timer]
OnCalendar=weekly

[Install]
WantedBy=basic.target
