# conf_server/systemd/system/restart-facebook.timer
[Unit]
Description=Run restart-facebook.service daily

[Timer]
OnCalendar=*-*-* 7:40:00
