# conf_server/systemd/system/transmission-stop.timer
[Unit]
Description=Run transmission-stop.service daily

[Timer]
OnCalendar=*-*-* 23:00:00

[Install]
WantedBy=basic.target
