# conf_server/systemd/system/transmission-start.timer
[Unit]
Description=Run transmission-start.service daily

[Timer]
OnCalendar=*-*-* 7:40:00

[Install]
WantedBy=basic.target
